# supervisord.conf
[program:node]
command = node bot.js
directory = /home/dotcloud/current
